# scripts/com.kylejensen.operator-ingestion.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  TEMPLATE FILE - DO NOT EDIT MANUALLY

  This file is a template for the LaunchD agent configuration.
  When you run `scripts/setup_launchd.sh install`, it will generate
  the actual plist file with paths customized for your system.

  The generated file will be placed at:
    ~/Library/LaunchAgents/com.kylejensen.operator-ingestion.plist

  To customize paths, set these before running setup_launchd.sh:
    export OPERATOR_DATA_DIR="/path/to/your/transcripts"

  Example paths that will be auto-detected:
    - Project root: /path/to/operator (detected from script location)
    - Script to run: ${PROJECT_ROOT}/scripts/daily_ingestion.sh
    - Working directory: ${PROJECT_ROOT}
    - Logs: ${PROJECT_ROOT}/ledger/logs/
-->
<plist version="1.0">
<dict>
    <!-- Job Label (must be unique) -->
    <key>Label</key>
    <string>com.kylejensen.operator-ingestion</string>

    <!-- Program to run (auto-generated with your project path) -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>${PROJECT_ROOT}/scripts/daily_ingestion.sh</string>
    </array>

    <!-- Working directory (auto-generated with your project path) -->
    <key>WorkingDirectory</key>
    <string>${PROJECT_ROOT}</string>

    <!-- Run every Monday at 9:00 AM -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Standard output log (auto-generated with your project path) -->
    <key>StandardOutPath</key>
    <string>${PROJECT_ROOT}/ledger/logs/launchd_stdout.log</string>

    <!-- Standard error log (auto-generated with your project path) -->
    <key>StandardErrorPath</key>
    <string>${PROJECT_ROOT}/ledger/logs/launchd_stderr.log</string>

    <!-- Don't run on load -->
    <key>RunAtLoad</key>
    <false/>

    <!-- Environment variables (OPERATOR_DATA_DIR from your shell) -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>OPERATOR_DATA_DIR</key>
        <string>${OPERATOR_DATA_DIR}</string>
    </dict>

    <!-- Don't keep alive -->
    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
