# scripts/com.operator.ledger-audit.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.operator.ledger-audit</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/kylejensen/Desktop/operator/scripts/ledger-audit.sh</string>
    </array>
    <key>StandardOutPath</key>
    <string>/Users/kylejensen/Desktop/operator/ledger/logs/audit-latest.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/kylejensen/Desktop/operator/ledger/logs/audit-latest.log</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
</dict>
</plist>
